# ingest-deno.service
[Unit]
Description=Service to forward sensor data to web
StartLimitIntervalSec=0

[Service]
Type=simple
Restart=always
RestartSec=1
WorkingDirectory=/home/b/bhong/ingest-deno2/
User=bhong
ExecStart=deno run --unstable  --allow-net --allow-read src/server.ts

[Install]
WantedBy=multi-user.target