# ingest-deno.service
[Unit]

Description=Service to forward sensor data to web

StartLimitIntervalSec=0



[Service]

Type=simple

Restart=always

RestartSec=1

User=root

ExecStart=/home/b/bhong/.deno/bin/deno run --unstable  --allow-net --allow-read /home/b/bhong/ingest-deno2/server.ts



[Install]

WantedBy=multi-user.target
