# ingest-deno.service
[Unit]
Description=Service to forward sensor data to web
StartLimitIntervalSec=0

[Service]
Type=simple
Restart=always
RestartSec=1
WorkingDirectory=/ingest-deno
ExecStart=/bin/sh -c '/ingest-deno/deno-bin/deno run --unstable --allow-net --allow-write=sensor-data.db,sensor-data.db-journal --allow-read --allow-env --allow-run --allow-sys --watch /ingest-deno/src/server.ts >> /ingest-deno/logs 2>&1'

[Install]
WantedBy=multi-user.target